# .config/systemd/user/bspwm-session.target
[Unit]
Description=Bspwm running
BindsTo=xsession.target
