# .config/runit/services/picom/run
#!/bin/bash

exec 2>&1

echo "Picom started" 

exec picom --experimental-backends --backend glx --xrender-sync-fence
