# .config/runit/services/mpd/run
#!/bin/bash

exec 2>&1

echo "MPD started" 

exec mpd --no-daemon
