# .config/runit/services/mpDris2/run
#!/bin/bash

exec 2>&1

# Read config file
test -x ./config && . ./config

echo "$SERVICE_NAME started" 

exec mpDris2
