# .config/runit/services/sxhkd/run
#!/bin/bash

exec 2>&1

echo "Sxhkd started" 

cd /home/sidharta

exec sxhkd
