# .config/runit/services/polybar/run
#!/bin/bash

exec 2>&1

echo "Polybar started" 

exec polybar
